# data/systemd/fragoi-autobright.service
[Unit]
Description=Automatic Brightness

[Service]
Type=exec
ExecStart=autobright

[Install]
WantedBy=graphical-session.target
